# systemd/parking-daemon-gate-out@.service
[Unit]
Description=E-Parking Gate-Out Daemon for Gate %i
After=network.target redis.service

[Service]
Type=simple
User=parking
Group=parking
SupplementaryGroups=dialout
WorkingDirectory=/opt/parking-system-v2
Environment=APP_ENV=production
Environment=PATH=/opt/parking-system-v2/.venv/bin
EnvironmentFile=/opt/parking-system-v2/.env
ExecStart=/opt/parking-system-v2/.venv/bin/python -m daemons.gate_out --gate-id %i
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=parking-daemon-gate-out-%i

[Install]
WantedBy=multi-user.target
