# systemd/parking-backup.service
[Unit]
Description=E-Parking Database Backup
After=postgresql.service

[Service]
Type=oneshot
User=parking
ExecStart=/opt/parking-system-v2/scripts/backup_database.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=parking-backup
